# app-admin/aws-cfn-bootstrap/files/cfn-hup.init
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	use dns logger
	need localmount net
}

start() {
	ebegin "Starting cfn-hup"
	start-stop-daemon --start --exec /usr/bin/cfn-hup
	eend $?
}

stop() {
	ebegin "Stopping cfn-hup"
	start-stop-daemon --stop --pidfile /var/run/cfn-hup.pid
	eend $?
}
